# camera/systemd/beadz-capture.service
[Unit]
Description=Beadz Fault-Cam 01 capture-and-sign (one-shot)
Wants=network-online.target

[Service]
Type=oneshot
User=beadz
Group=beadz
SupplementaryGroups=video
EnvironmentFile=/etc/beadz-camera/device.env
WorkingDirectory=/opt/beadz-camera
ExecStart=/opt/beadz-camera/venv/bin/beadz-camera --env /etc/beadz-camera/device.env capture-once
# PrivateTmp isolates /tmp; bring-up (smoke.sh step 0) verifies /tmp is tmpfs so raw frames never touch the SD card
PrivateTmp=yes
NoNewPrivileges=yes
