# camera/systemd/beadz-capture.service
[Unit]
Description=Beadz Fault-Cam 01 capture-and-sign (one-shot)
Wants=network-online.target
After=tmp.mount

[Service]
Type=oneshot
User=beadz
Group=beadz
SupplementaryGroups=video
EnvironmentFile=/etc/beadz-camera/device.env
WorkingDirectory=/opt/beadz-camera
# Fail-closed: skip the capture (logged as condition-failed) rather than ever
# write a raw frame to an SD-backed /tmp. smoke.sh step 0 is the bring-up
# check; this is the per-run authority.
ExecCondition=/usr/bin/findmnt -n -t tmpfs /tmp
ExecStart=/opt/beadz-camera/venv/bin/beadz-camera --env /etc/beadz-camera/device.env capture-once
# PrivateTmp isolates /tmp; bring-up (smoke.sh step 0) verifies /tmp is tmpfs so raw frames never touch the SD card
PrivateTmp=yes
NoNewPrivileges=yes
